# systemd/options_ai_storage_monitor.service
[Unit]
Description=OptionsPredicator Storage Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nate
Group=nate
WorkingDirectory=/opt/OptionsPredicator
EnvironmentFile=/opt/OptionsPredicator/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/OptionsPredicator/.venv/bin/python -m options_ai.storage_monitor_main
Restart=always
RestartSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
