# systemd/spx_chain_phase2_term_dte21t3.service
# Phase 2 TERM: DTE21±3 chain feature/label builder (TimescaleDB)
#
# Install:
#   sudo cp systemd/spx_chain_phase2_term_dte21t3.service /etc/systemd/system/spx_chain_phase2_term_dte21t3.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now spx_chain_phase2_term_dte21t3
#
# Logs:
#   journalctl -u spx_chain_phase2_term_dte21t3 -f

[Unit]
Description=SPX TERM chain feature/label builder (dte21t3)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nate
Group=nate
WorkingDirectory=/opt/OptionsPredicator

EnvironmentFile=/opt/OptionsPredicator/.env
Environment=PYTHONUNBUFFERED=1

Environment=TARGET_DTE_DAYS=21
Environment=DTE_TOLERANCE_DAYS=3
Environment=TERM_BUCKET=term_dte21t3

# half/full horizons (calendar days): 11d, 21d
Environment=HORIZONS_MINUTES=15840,30240
Environment=MAX_FUTURE_LOOKAHEAD_MINUTES=20160

ExecStart=/opt/OptionsPredicator/.venv/bin/python -m options_ai.phase2_term_main
Restart=always
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
